# mia_jetson.service
[Unit]
Description=MIA JETSON AI Assistant
After=network.target sound.target
Wants=sound.target

[Service]
ExecStart=/usr/bin/python3 -u src/main.py
WorkingDirectory=/home/mia/MIA_JETSON
User=mia
SupplementaryGroups=audio

# --- Restart policy ---
Restart=always
RestartSec=5
StartLimitBurst=10
StartLimitInterval=60s

# --- Logs ---
StandardOutput=append:/home/mia/MIA_JETSON/logs/mia.log
StandardError=append:/home/mia/MIA_JETSON/logs/mia.log

# --- Environment ---
# Give ALSA/PulseAudio time to settle after boot before starting
ExecStartPre=/bin/sleep 5

[Install]
WantedBy=multi-user.target
